# Library/LaunchAgents/me.ilynxcat.gitd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>me.ilynxcat.gitd</string>
		<key>ProgramArguments</key>
		<array>
			<string>/opt/homebrew/bin/git</string>
			<string>daemon</string>
			<string>--reuseaddr</string>
			<string>--listen=127.0.0.1</string>
			<string>--base-path=/usr/local/repos</string>
			<string>--enable=receive-pack</string>
			<string>--export-all</string>
			<string>/usr/local/repos</string>
		</array>
		<key>RunAtLoad</key>
		<true />
		<key>KeepAlive</key>
		<true />
	</dict>
</plist>
